# misc/netpoke.service
[Unit]
Description=NetPoke WiFi Verification Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netpoke
ExecStart=/usr/local/bin/netpoke-server
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
